# bin/dht22-server.service
[Unit]
Description=DHT22 Server Spring Boot Application
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/dht22-server
ExecStart=/usr/bin/java -Xms32m -Xmx128m -XX:+UseSerialGC -jar /home/%i/dht22-server/dht22-server.jar
StandardOutput=null
StandardError=null
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target